# client/boot-tools/includes/etc/profile.d/loadenviron.sh
#!/bin/bash

### NOTA este archivo se debe de copiar desde el fichero etc/init/default

